# clearing-house-app/clearing-house-api/supervisord.conf
[supervisord]
nodaemon=true

[program:clearing-house-api]
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
command=/server/clearing-house-api
autorestart=true
